# source/upstart/nginx.conf
start on startup
respawn
exec /usr/sbin/nginx
