# source/upstart/vndf-game-service.conf
description "Von Neumann Defense Force daemon"

start on startup
respawn

chdir /opt
exec /opt/vndf-game-service
