# source/upstart/vndf-core-service.conf
description "Von Neumann Defense Force daemon"

start on startup
respawn
exec /opt/vndf/vndf-core-service
