# source/upstart/vndf-server.conf
description "Main Von Neumann Defense Force daemon"

start on startup
exec /opt/vndf/vndf-server
